# hwexp/packaging/hwexp.service
[Unit]
Description=AUDiot Hardware Telemetry Exporter
After=network.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/bin/hwexp --config /etc/hwexp/hwexp.yaml --fixture /etc/hwexp/sample_hwmon.json
Restart=on-failure
RestartSec=5s

# Security hardening
ProtectSystem=full
PrivateTmp=true

[Install]
WantedBy=multi-user.target
